# lander/etc/xinetd.d/lander
service lander
{
  type = UNLISTED
  instances = 10
  socket_type = stream
  protocol = tcp
  wait = no
  user = root
  port = 9999
  server = /usr/local/bin/lander-handler.sh
}
